# deploy/llmgateway-backup.timer
[Unit]
Description=Run LLMGateway encrypted backup every six hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:15:00
Persistent=true
RandomizedDelaySec=15m
Unit=llmgateway-backup.service

[Install]
WantedBy=timers.target
