# deploy/llmgateway-backup.service
[Unit]
Description=LLMGateway encrypted remote S3 backup
Requires=docker.service
Wants=network-online.target
After=docker.service network-online.target llmgateway-compose.service
RequiresMountsFor=/opt/llmgateway /etc/llmgateway /etc/llmgateway-backup /var/lib/llmgateway-backup

[Service]
Type=oneshot
ExecStart=/opt/llmgateway/backup-bundle-launcher-linux.sh scheduled-backup /etc/llmgateway-backup/backup.env
TimeoutStartSec=85m
User=root
Group=root
UMask=0077
RuntimeDirectory=llmgateway-restic
RuntimeDirectoryMode=0700
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/llmgateway /etc/llmgateway-backup /opt/llmgateway
ReadWritePaths=/var/lib/llmgateway-backup /run/lock /run/llmgateway-restic
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
